# launchagents/com.moba_ranking.post_to_x_unite_tier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- 一意な識別子 -->
    <key>Label</key>
    <string>com.moba_ranking.post_to_x_unite_tier</string>
    
    <!-- 実行するプログラムと引数 -->
    <key>ProgramArguments</key>
    <array>
      <!-- pythonのパス（仮想環境） -->
      <string>/Users/yamamotokazuki/develop/moba-ranking/moba-ranking-venv/bin/python</string>
      <!-- スクリプトの相対パス -->
      <string>post_to_x_unite_tier.py</string>
    </array>
    
    <!-- 作業ディレクトリの設定 -->
    <key>WorkingDirectory</key>
    <string>/Users/yamamotokazuki/develop/moba-ranking</string>
    
    <!-- 標準出力と標準エラーのログ出力先 -->
    <key>StandardOutPath</key>
    <string>/Users/yamamotokazuki/develop/moba-ranking/logs/com.moba_ranking.post_to_x_unite_tier.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yamamotokazuki/develop/moba-ranking/logs/com.moba_ranking.post_to_x_unite_tier.log</string>
    
    <!-- 毎週金曜日の18時00分に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>00</integer>
    </dict>
  </dict>
</plist>
